# supervisor/etc/supervisord_d/helloworld-prod-go.ini
[program:helloworld-prod-go]
command=/path/to/helloworld-prod-go/app/main
autostart = true
startsecs = 5
user = root
redirect_stderr = true
stdout_logfile = /var/log/supervisor/helloworld-prod-go.log
